# sv/swarm-client/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

exec 2>&1
exec java -jar /swarm-client.jar -master $MASTER -deleteExistingClients -disableClientsUniqueId
